# .crontab
# Mac preferences backup
@reboot /bin/sleep 300; /usr/bin/rsync -avP --delete --exclude=.DS_Store /Users/gabrielecanepa/Library/Preferences/ /Users/gabrielecanepa/Library/Mobile\ Documents/com~apple~CloudDocs/System/Mac\ Preferences >> /var/log/crontab.log 2>&1
